# upstart.conf
start on runlevel [234]
respawn
env CANTAL_LOG=warn
exec /usr/bin/cantal-agent --host 127.0.0.1 --port 22682

